# init.d/auditd.cron
#!/bin/sh

##########
# This script can be installed to the cron system to get log rotation
# based on time instead of log size.
##########

/sbin/auditctl --signal usr1
EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
    /usr/bin/logger -t auditd "ALERT auditctl exited abnormally with [$EXITVALUE] while rotating the logs"
fi
exit 0

